# data/s6-overlay/s6-rc.d/init-mailer-config/up
/etc/s6-overlay/s6-rc.d/init-mailer-config/run
