# data/s6-overlay/s6-rc.d/svc-opendkim/run
#!/usr/bin/with-contenv bash

/usr/sbin/opendkim -D -f -x /etc/opendkim/opendkim.conf
