# data/s6-overlay/s6-rc.d/svc-postfix/run
#!/usr/bin/with-contenv bash

/usr/sbin/postfix -c /etc/postfix start-fg
